# hooks/systemd-update-efi.path
[Unit]
Description=Update systemd-boot

[Path]
PathChanged=/usr/bin/bootctl

[Install]
WantedBy=multi-user.target
WantedBy=system-update.target
