# hooks/systemd-update-efi.service
[Unit]
Description=Update systemd-boot

[Service]
Type=oneshot
ExecStart=/usr/bin/bootctl update
